# virttmp/rc.local
#!/bin/sh
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

/bin/bash /bin/setipbymac-ubuntu.sh H1-ubuntu14.04 192.168.122.1 0.0.0.0 0.0.0.0 52:54:00:C3:49:4F 10.0.20.30 255.255.255.0 52:54:00:07:D6:A2
