# virttmp/rc.local
#!/bin/sh
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

/bin/bash /bin/setipbymac-ubuntu.sh ubuntu14.04-resize 192.168.122.1 0.0.0.0 0.0.0.0 52:54:00:92:CC:BE 10.0.20.30 255.255.255.0 52:54:00:94:89:A7
