# virttmp/rc.local
#!/bin/sh
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

/bin/bash /bin/setipbymac-ubuntu.sh ubuntu14.04-test2 192.168.122.1 0.0.0.0 0.0.0.0 52:54:00:10:9D:07 10.0.20.20 255.255.255.0 52:54:00:96:CA:C8
